# CRM/deploy/edugrowth-crm.service
[Unit]
Description=EduGrowth CRM API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/edugrowth-crm
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start
Restart=always
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
